# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "pid /tmp/nginx.pid; daemon off;"

[program:nexus]
command=/opt/sonatype/nexus/bin/nexus run
